# dashboard/deploy/dashboard-backend.service
[Unit]
Description=Overseer Ops Backend (orphan detection / git sync / project progress)
After=network.target

[Service]
User=root
WorkingDirectory=/mnt/ssd/overseeer/dashboard/server
ExecStart=/mnt/ssd/overseeer/dashboard/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8010
Restart=always
RestartSec=5
Environment=PYTHONPATH=/mnt/ssd/overseeer/dashboard/server

[Install]
WantedBy=multi-user.target
